# templates/systemd/spacy-release-api.service
[Unit]
Description=SpacyRelease Feedback API (Port 5001)
After=network.target mariadb.service
Requires=mariadb.service

[Service]
Type=simple
User=spacy-filter
Group=spacy-filter
WorkingDirectory=/opt/spacyserver/api
Environment=PYTHONPATH=/opt/spacyserver:/opt/spacyserver/api
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/spacyserver/venv/bin/python3 /opt/spacyserver/api/quarantine_feedback.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=spacy-release-api

# Resource limits
MemoryMax=256M
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
