# templates/systemd/spacy-whitelist-api.service
[Unit]
Description=SpacyWhitelist API (Port 5002)
After=network.target

[Service]
Type=simple
User=spacy-filter
Group=spacy-filter
WorkingDirectory=/opt/spacyserver/api
Environment=PYTHONPATH=/opt/spacyserver:/opt/spacyserver/api
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/spacyserver/venv/bin/python3 /opt/spacyserver/api/whitelist_notification.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=spacy-whitelist-api

# Resource limits
MemoryMax=256M
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
